# debian/postinst
#!/bin/bash

#DEBHELPER#

# Install plugins
echo Installing Grafana plugins
sudo grafana-cli plugins install ae3e-plotly-panel

# Enable Grafana web UI port
echo "Enabling Grafana TCP port 3000"
sudo ufw allow 3000/tcp

# Start Grafana service
echo "Starting Grafana server service"
sudo /bin/systemctl start grafana-server

# Enable Grafana service
echo "Enabling Grafana server service"
sudo /bin/systemctl enable grafana-server

# Fixed sleep is required, otherwise installation fails during image build process
echo "Waiting for 15 seconds for Grafana service to start"
sleep 15

# Create service account
echo "Creating Grafana service account"
curl -s -X POST http://admin:admin@127.0.0.1:3000/api/serviceaccounts -H "Accept: application/json" -H "Content-Type: application/json" -d '{"name": "wlanpi", "role": "Admin", "isDisabled": false}'
echo

# Generate service API token
echo "Generating Grafana API token"
GRAFANA_TOKEN=$(curl -s -X POST http://admin:admin@127.0.0.1:3000/api/serviceaccounts/2/tokens -H "Accept: application/json" -H "Content-Type: application/json" -d '{"name": "wlanpi", "role": "Admin"}' | jq -r '.key')

# Save API token to env variable
echo "Saving Grafana API token to env variable"
echo "GRAFANA_TOKEN=$GRAFANA_TOKEN" | sudo tee -a /etc/environment
set -a; source /etc/environment; set +a