# debian/postinst
#!/bin/bash

#DEBHELPER#

# Enable Grafana web UI port
echo "Enabling TCP port 3000"
sudo ufw allow 3000/tcp

# Enable Grafana service
echo "Enabling grafana-server service"
sudo /bin/systemctl enable grafana-server

# Start Grafana service
echo "Enabling grafana-server service"
sudo /bin/systemctl start grafana-server