# debian/postinst
#!/bin/bash

#DEBHELPER#

# Install plugins
echo Installing grafana plugins
sudo grafana-cli plugins install ae3e-plotly-panel

# Enable Grafana web UI port
echo "Enabling TCP port 3000"
sudo ufw allow 3000/tcp

# Enable Grafana service
echo "Enabling grafana-server service"
sudo /bin/systemctl enable grafana-server

# Start Grafana service
echo "Enabling grafana-server service"
sudo /bin/systemctl start grafana-server

# Create service account
curl -s -X POST http://admin:admin@localhost:3000/api/serviceaccounts -H "Accept: application/json" -H "Content-Type: application/json" -d '{"name": "wlanpi", "role": "Admin", "isDisabled": false}'

# Generate service API token
GRAFANA_TOKEN=$(curl -s -X POST http://admin:admin@localhost:3000/api/serviceaccounts/2/tokens -H "Accept: application/json" -H "Content-Type: application/json" -d '{"name": "wlanpi", "role": "Admin"}' | jq -r '.key')

# Save API token to env variable
echo "GRAFANA_TOKEN=$GRAFANA_TOKEN" | sudo tee -a /etc/environment
set -a; source /etc/environment; set +a