# debian/postinst
#!/bin/bash

#DEBHELPER#

# Install plugins
echo Installing Grafana plugins
sudo grafana-cli plugins install ae3e-plotly-panel

# Start Grafana service
echo "Starting Grafana server service"
sudo /bin/systemctl stop grafana-server
sudo /bin/systemctl start grafana-server

# Enable Grafana service
echo "Enabling Grafana server service"
sudo /bin/systemctl enable grafana-server
